# openwrt-imagebuilder/overlay_old/etc/uci-defaults/96-landingpage.sh
# -- Set up Landing Page configuration

touch /etc/config/landingpage
uci add landingpage settings
uci set landingpage.@settings[0].autostart='0'

uci commit landingpage

exit 0