# overlay/etc/init.d/resetssids
#!/bin/sh /etc/rc.common
# Copyright (c) 2019 Hak5 LLC

START=90

start() {
	/usr/sbin/resetssids &>/dev/null &
	sleep 1
}

boot() {
	start
}